# util/systemd/hive-queen.service
[Unit]
Description=Queen Launcher
After=couchdb.service openntpd.service systemd-remount-fs.service
Requires=couchdb.service

[Service]
Type=simple
Restart=always
RestartSec=15s
ExecStart=/bin/node /root/GroundHive/Queen/server.js

[Install]
WantedBy=multi-user.target


