# util/systemd/hive-queen.service
[Unit]
Description=Hive Queen Launcher
After=couchdb.service
Requires=couchdb.service

[Service]
Type=simple
Restart=always
RestartSec=10s
ExecStart=/usr/bin/node /root/Hive/Queen/server.js

[Install]
WantedBy=multi-user.target
