# util/systemd/hive-honeycomb.service
[Unit]
Description=Hive Honeycomb Launcher
After=couchdb.service
Requires=couchdb.service

[Service]
Type=simple
Restart=always
RestartSec=10s
ExecStart=/usr/bin/node /root/Hive/Honeycomb/server.js

[Install]
WantedBy=multi-user.target
