# util/systemd/hive-router.service
[Unit]
Description=Hive Router Launcher
After=couchdb.service
Requires=couchdb.service

[Service]
Type=simple
Restart=always
RestartSec=10s
ExecStart=/usr/bin/node /root/Hive/Router/server.js

[Install]
WantedBy=multi-user.target
